# crontab
*/1 * * * * aplay /docker-speaker-standby/10.wav >> /var/log/speaker_standby.log 2>&1
