# crontab
*/5 * * * * aplay /docker-speaker-standby/10.wav

