# crontab
*/5 * * * * aplay /docker-speaker_standby/10.wav

